# dot_config/systemd/user/waybar.service.d/override.conf
[Service]
RestartSec=3sec
Slice=background.slice
